# examples/httpreqd/systemd/daemonpp.service.in
# Properties docs: https://www.freedesktop.org/software/systemd/man/systemd.service.html
[Unit]
Description=@PROJECT_DESCRIPTION@
After=network.target

[Service]
Type=forking
ExecStart=/usr/bin/@PROJECT_NAME@ --config /etc/@PROJECT_NAME@/@PROJECT_NAME@.conf
ExecReload=/bin/kill -s SIGHUP $MAINPID
ExecStop=/bin/kill -s SIGTERM $MAINPID
User=root
StandardError=syslog
SyslogIdentifier=@PROJECT_NAME@

[Install]
WantedBy=multi-user.target